# service/tor/run
#!/bin/sh
if ! [ -f /etc/tor/torrc ]; then
    echo "Using torrc in /etc/service/tor/torrc"
    tor -f ./torrc
else
    echo "Using torrc in default location"
    tor -f /etc/tor/torrc
fi
